# config/cron.d/monitoring
* * * * * root /root/cron/monitor_and_restart_apache.py {{ public_hostname }} >> /root/cron/logfile 2>&1

